# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:bgutil]
# Servidor que gera PO Token automaticamente via BotGuard
# Fica rodando na porta 4416, o yt-dlp se conecta automaticamente
command=python -m bgutil_ytdlp_pot_provider server --host 127.0.0.1 --port 4416
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/bgutil.err.log
stdout_logfile=/var/log/supervisor/bgutil.out.log
# Aguarda 5s para o bgutil inicializar antes de subir a API
startsecs=5

[program:api]
command=gunicorn -b 0.0.0.0:5000 --timeout 180 --workers 2 app:app
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api.err.log
stdout_logfile=/var/log/supervisor/api.out.log
# Aguarda bgutil subir primeiro
startretries=3

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock
